# examples/user-consul-go/app/root/etc/services.d/app/run
#!/usr/bin/with-contenv sh

exec /app/main;